# phase6/systemd/llm-backup-daily.timer
[Unit]
Description=Run LLM backup once per day

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=llm-backup-daily.service

[Install]
WantedBy=timers.target
